# conf/cron/crontab
# m h  dom mon dow  user  command
*/5 * * * * TYPO3_CONTEXT=Development/Local /usr/local/bin/php /app/public/typo3/sysext/core/bin/typo3 scheduler:run
